# digitalocean/reboot.service
[Unit]
Description=Reboot

[Service]
Type=oneshot
ExecStart=/bin/sh -c "systemctl reboot"

[Install]
WantedBy=multi-user.target